# Install/Files/pulseaudio.service
[Unit]
Description=Sound Service
Before=naoqi.service

[Service]
Type=simple
ExecStart=/bin/true
Restart=on-failure
Slice=naoqi.slice
