# Install/Network/wired
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

BASE_ETH0_IP="192.168.%teamID%."
ROBOT_ID="%robotPart%"

depend() {
  need dbus
}

start() {
  # Replace buggy driver
  networkDriver=$(lspci -k -s 02:00.0 | grep "Kernel driver in use:" | awk 'NF>1{print $NF}')
  if [[ ${networkDriver} == "r8169" ]] ; then
    echo "0000:02:00.0" > /sys/bus/pci/drivers/r8169/unbind
    insmod /lib/modules/2.6.33.9-rt31-aldebaran-rt/kernel/drivers/net/r8168.ko
  fi

  ifconfig eth0 "${BASE_ETH0_IP}${ROBOT_ID}" netmask 255.255.0.0
  ifconfig eth0 up
}

stop() {
  ifconfig eth0 down
}

restart() {
  stop
  sleep 2
  start
}
